# roles/kube-proxy/templates/kube-proxy.service.j2
[Unit]
Description=Kubernetes Kube Proxy
Documentation=https://github.com/kubernetes/kubernetes

[Service]
ExecStart=/usr/local/bin/kube-proxy \
  --config=/etc/kubernetes/config/kube-proxy.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target